# Fedora/usb_lock_and_power_shutdown/80-power-shutdown.rules
ACTION=="change", SUBSYSTEM=="power_supply", ATTR{online}=="0", RUN+="/usr/local/bin/usb-lock_power-shutdown.sh shutdown"
